# supervisord.api.conf
; supervisord.api.conf — HF Spaces entrypoint (PID 1), API role
; Runs only FastAPI. Used when SPACE_ROLE=api (the default — see entrypoint.sh).
; The Celery worker runs in a separate HF Space (SPACE_ROLE=worker) so the
; two processes don't compete for RAM inside one small free-tier container.

[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid

[program:fastapi]
command=uvicorn app.main:app --host 0.0.0.0 --port 7860
directory=/code
autostart=true
autorestart=true
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
