# supervisord.conf
; supervisord.conf — HF Spaces entrypoint (PID 1), default/combined role
; Runs FastAPI (uvicorn) and the Celery ingestion worker as sibling processes
; in one container. This is the default (SPACE_ROLE unset) because Hugging
; Face's free tier only allows ONE non-static Space per account — a second
; Docker Space requires a PRO subscription, so the "split into two free
; Spaces" plan isn't actually free on HF and was reverted.
;
; supervisord.worker.conf (celery-only) is kept around in case a genuinely
; free second host shows up later (e.g. a VM elsewhere) — see entrypoint.sh.

[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid

[program:fastapi]
command=uvicorn app.main:app --host 0.0.0.0 --port 7860
directory=/code
autostart=true
autorestart=true
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:celery]
command=celery -A app.worker.celery_app worker --loglevel=info --concurrency=1
directory=/code
autostart=true
autorestart=true
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
