# supervisord.worker.conf
; supervisord.worker.conf — HF Spaces entrypoint (PID 1), worker role
; Runs only the Celery ingestion worker (+ a trivial health endpoint so HF's
; port-7860 health check passes). Used when SPACE_ROLE=worker — see
; entrypoint.sh. The FastAPI app runs in the main Space (SPACE_ROLE=api),
; so this container never imports app.main and never pays for that memory.

[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid

[program:healthcheck]
command=python -m app.worker.healthcheck
directory=/code
autostart=true
autorestart=true
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:celery]
command=celery -A app.worker.celery_app worker --loglevel=info --concurrency=1
directory=/code
autostart=true
autorestart=true
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
